# hwa/acpi/susp-v1/post-wake-usb-reset.service
[Unit]
Description=Reset USB devices after wake-up
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/post-wake-usb-reset.sh

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
